# sources/meta-ess/recipes-ess/services/files/host-setup.service
[Unit]
Description=Host init
After=network.target

[Service]
Type=oneshot
ExecStart=@SCRIPTDIR@/host-init.sh

[Install]
WantedBy=multi-user.target
