# sources/meta-ess/recipes-ess/services/files/ifc14xx-init.service
[Unit]
Description=IFC14XX init
After=systemd-modules-load.service

[Service]
Type=oneshot
ExecStart=@SCRIPTDIR@/ifc14xx-init.sh

[Install]
WantedBy=multi-user.target
